# tmux-backup/2026-02-13/systemd/tmux-zai-key-bootstrap.service
[Unit]
Description=Bootstrap ZAI key into tmux web session
After=ttyd-tmux.service
Wants=ttyd-tmux.service

[Service]
Type=oneshot
Environment=TMUX_BOOT_SESSION=web
Environment=TMUX_BOOT_WAIT_SEC=45
ExecStart=/root/.local/bin/tmux-boot-key-sync.sh

[Install]
WantedBy=multi-user.target
