# tmux-backup/2026-02-13/systemd/ttyd-tmux.service
[Unit]
Description=ttyd web terminal (tmux backend)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=~
Environment=LANG=ko_KR.UTF-8 LC_ALL=ko_KR.UTF-8
KillMode=process
ExecStartPre=/bin/bash -lc '/usr/bin/tmux has-session -t web 2>/dev/null || /usr/bin/tmux new-session -d -s web'
ExecStart=/usr/bin/ttyd -p 7681 -i lo -W -t disableResizeOverlay=true -t cursorStyle=bar /usr/bin/tmux attach-session -t web
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
